# systemd/spbu-tt-diff-notify.timer
[Unit]
Description=Run Spbu Timetable changes notifier

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
